# supervisord.conf
[program:node]
command = node lib/server.js
directory = /home/dotcloud/current